# plugins/deploy/templates/systemd.service
[Unit]
Description={{DOMAIN}}
After=network.target

[Service]
Type=simple
User=deployer
WorkingDirectory={{DIR}}/release
# {{RUBY}} dir must be on PATH so `bundle` resolves under -S.
Environment="PATH={{RUBY_DIR}}:/home/deployer/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart={{RUBY}} -S bundle exec lux s -p {{PORT}} -e production
Restart=always
RestartSec=5
EnvironmentFile=-{{DIR}}/.env

[Install]
WantedBy=multi-user.target
